# docker/php/cron/crontab
* * * * * root /usr/local/bin/php /www/yii schedule/run --scheduleFile=@console/config/schedule.php >> /var/log/cron.log 2>&1
